# configs/rpi2/rootfs_cfg/etc/crontab
0 * * * * ntpdate  ntp.univ-lyon1.fr canon.inria.fr ntp1.tuxfamily.net 
