# configs/rpi3/rootfs_cfg/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:-/bin/sh

::ctrlaltdel:/sbin/reboot


